# .config/systemd/user/battery-alert.service
[Unit]
Description=Battery alert service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-alert

[Install]
WantedBy=default.target
